# Linux/Units/ssh-host-keygen.service/ssh-host-keygen.service
[Unit]
Description=Generate SSH Host Keys if missing
Before=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
